# rootfs/etc/services.d/apt-packages-on-run/up
#!/usr/bin/execlineb -P
apt-get-install-min $packages
